# case-ansible/roles/app/templates/ty-case.ini
[program:ty-case]
command=/var/www/.env/bin/python /var/www/ty-case/main.py
numprocs=1
numprocs_start=0
priority=999
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=1
directory=/var/www/ty-case/
user=root
redirect_stderr=false
stdout_logfile=/var/log/supervisor/ty-case-out.log
stderr_logfile=/var/log/supervisor/ty-case-errs.log